# modules/microservices/templates/static_assets.service.erb
[Unit]
Description=NewsFeed Quotes Service

[Service]
Type=simple


User=root
SuccessExitStatus=143
ExecStart=cd /usr/share/infra-problem/front-end/public;/usr/bin/python3.4 /usr/share/infra-problem/front-end/public/serve.py

[Install]
WantedBy=default.target
